# tribus/config/data/tribus.supervisor.conf
[program:tribus-celery]
command=/usr/bin/python /media/desarrollo/tribus/manage.py celeryd
directory=/media/desarrollo/tribus
user=www-data
numprocs=1
stdout_logfile=/var/log/tribus/celeryd.log
stderr_logfile=/var/log/tribus/celeryd.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30

[program:tribus-celerybeat]
command=/usr/bin/python /media/desarrollo/tribus/manage.py celerybeat
directory=/media/desarrollo/tribus
user=www-data
numprocs=1
stdout_logfile=/var/log/tribus/celerybeat.log
stderr_logfile=/var/log/tribus/celerybeat.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30

